# nethunter-fs/kali-armhf/var/lib/dpkg/info/john.list
/.
/etc
/etc/john
/etc/john/john-mail.msg
/etc/john/john.conf
/etc/john/john-mail.conf
/etc/cron.d
/etc/cron.d/john
/usr
/usr/sbin
/usr/sbin/john
/usr/sbin/mailer
/usr/share
/usr/share/john
/usr/share/john/dynamic.conf
/usr/share/john/repeats32.conf
/usr/share/john/regex_alphabets.conf
/usr/share/john/dumb16.conf
/usr/share/john/dumb32.conf
/usr/share/john/korelogic.conf
/usr/share/john/cronjob
/usr/share/john/repeats16.conf
/usr/share/john/dynamic_flat_sse_formats.conf
/usr/share/john/john.local.conf
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/unique.8.gz
/usr/share/man/man8/mailer.8.gz
/usr/share/man/man8/unshadow.8.gz
/usr/share/man/man8/john.8.gz
/usr/share/man/man8/unafs.8.gz
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/john.bash_completion
/usr/share/doc
/usr/share/doc/john
/usr/share/doc/john/EXTERNAL.gz
/usr/share/doc/john/CREDITS
/usr/share/doc/john/EXAMPLES.gz
/usr/share/doc/john/README.gz
/usr/share/doc/john/RULES.gz
/usr/share/doc/john/MODES.gz
/usr/share/doc/john/changelog.Debian.gz
/usr/share/doc/john/copyright
/usr/share/doc/john/CONFIG.gz
/usr/share/doc/john/changelog.gz
/usr/share/doc/john/FAQ.gz
/usr/share/doc/john/OPTIONS.gz
/usr/share/doc/john/README.Debian
/usr/sbin/unafs
/usr/sbin/unique
/usr/sbin/unshadow
