# nethunter-fs/kali-armhf/var/lib/dpkg/info/tinyproxy.conffiles
/etc/init.d/tinyproxy
/etc/logrotate.d/tinyproxy
/etc/tinyproxy.conf
