# nethunter-fs/kali-armhf/var/lib/dpkg/info/tinyproxy.prerm
#!/bin/sh
set -e
# Automatically added by dh_installinit
if [ -x "/etc/init.d/tinyproxy" ]; then
	invoke-rc.d tinyproxy stop || exit $?
fi
# End automatically added section
