# nethunter-fs/kali-armhf/var/lib/dpkg/info/tinyproxy.list
/.
/usr
/usr/share
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/tinyproxy
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/tinyproxy.8.gz
/usr/share/man/man5
/usr/share/man/man5/tinyproxy.conf.5.gz
/usr/share/tinyproxy
/usr/share/tinyproxy/default.html
/usr/share/tinyproxy/debug.html
/usr/share/tinyproxy/stats.html
/usr/share/doc
/usr/share/doc/tinyproxy
/usr/share/doc/tinyproxy/changelog.Debian.armhf.gz
/usr/share/doc/tinyproxy/filter-howto.txt
/usr/share/doc/tinyproxy/README.Debian
/usr/share/doc/tinyproxy/http-error-codes.txt
/usr/share/doc/tinyproxy/NEWS.gz
/usr/share/doc/tinyproxy/changelog.Debian.gz
/usr/share/doc/tinyproxy/changelog.gz
/usr/share/doc/tinyproxy/copyright
/usr/share/doc/tinyproxy/http-rfcs.txt
/usr/share/doc/tinyproxy/examples
/usr/share/doc/tinyproxy/examples/tinyproxy.conf.gz
/usr/sbin
/usr/sbin/tinyproxy
/etc
/etc/init.d
/etc/init.d/tinyproxy
/etc/logrotate.d
/etc/logrotate.d/tinyproxy
/etc/tinyproxy.conf
/var
/var/log
/var/log/tinyproxy
